# vm/supervisord.ini
---------------------- AutoAdd + AutoUpdate | 45.32.62.170 -------------------------------------------------------------

[program:character_add]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_add
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_friends_add]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_friends_add
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_achievements_add]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_achievements_add
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:free_company_add]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest free_company_add
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:linkshell_add]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest linkshell_add
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:pvp_team_add]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest pvp_team_add
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_update]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_update
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_friends_update]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_friends_update
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_achievements_update]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_achievements_update
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:free_company_update]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest free_company_update
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:linkshell_update]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest linkshell_update
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:pvp_team_update]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest pvp_team_update
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

---------------------- AutoUpdate (6 consumers) | 202.182.123.22 -------------------------------------------------------

[program:character_update_0_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_update_0_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_update_1_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_update_1_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_update_2_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_update_2_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_update_3_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_update_3_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_update_4_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_update_4_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_update_5_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_update_5_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

---------------------- AutoUpdateFriends | 207.148.110.192 -------------------------------------------------------------

[program:character_friends_update_0_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_friends_update_0_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_friends_update_1_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_friends_update_1_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_friends_update_0_patreon]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_friends_update_0_patreon
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_friends_update_1_patreon]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_friends_update_1_patreon
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

---------------------- AutoUpdatePvpTeam | 198.13.34.9 -----------------------------------------------------------------

[program:pvp_team_update_0_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest pvp_team_update_0_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:pvp_team_update_1_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest pvp_team_update_1_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:pvp_team_update_0_patron]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest pvp_team_update_0_patron
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:pvp_team_update_1_patron]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest pvp_team_update_1_patron
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

---------------------- Patreon | 198.13.61.206 -------------------------------------------------------------------------
--- I should move all patron queues to this server

[program:character_update_0_patreon]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_update_0_patreon
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_update_1_patreon]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_update_1_patreon
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

---------------------- AutoUpdateFreeCompany | 45.76.199.186 -----------------------------------------------------------

[program:free_company_update_0_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest free_company_update_0_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:free_company_update_1_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest free_company_update_1_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:free_company_update_0_patron]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest free_company_update_0_patron
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:free_company_update_1_patron]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest free_company_update_1_patron
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

---------------------- AutoUpdateLinkshells | 108.160.137.59 -----------------------------------------------------------

[program:linkshell_update_0_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest linkshell_update_0_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:linkshell_update_1_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest linkshell_update_1_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:linkshell_update_0_patron]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest linkshell_update_0_patron
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:linkshell_update_1_patron]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest linkshell_update_1_patron
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

---------------------- AutoUpdateAchievements | 45.76.107.66 -----------------------------------------------------------
-- todo: add another queue for this, move patrons to their own server


[program:character_achievements_update_0_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_achievements_update_0_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_achievements_update_1_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_achievements_update_1_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_achievements_update_0_patreon]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_achievements_update_0_patreon
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_achievements_update_1_patreon]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_achievements_update_1_patreon
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log
