# vm/supervisord.ini
---------------------- Server 1 - Adding ----------------------------------------------------------------

[program:character_add]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_add
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_friends_add]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_friends_add
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_achievements_add]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_achievements_add
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:free_company_add]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest free_company_add
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:linkshell_add]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest linkshell_add
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:pvp_team_add]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest pvp_team_add
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

---------------------- Server 2 - Manual Updating ----------------------------------------------------------------

[program:character_update]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_update
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_friends_update]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_friends_update
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_achievements_update]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_achievements_update
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:free_company_update]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest free_company_update
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:linkshell_update]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest linkshell_update
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:pvp_team_update]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest pvp_team_update
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

---------------------- Server 3 - Auto Updating 0,1,2,3 ----------------------------------------------------------------

[program:character_update_0_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_update_0_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_update_1_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_update_1_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_update_2_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_update_2_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_update_3_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_update_3_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

---------------------- Server 4 - Auto Updating 4,5 ----------------------------------------------------------------

[program:character_update_4_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_update_4_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_update_5_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_update_5_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

---------------------- Server 5 - Auto Updating - Low Priority ----------------------------------------------------------------

[program:character_update_0_low]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_update_0_low
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_update_1_low]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_update_1_low
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

---------------------- Server 6 - Auto Updating - Patreons ----------------------------------------------------------------

[program:character_update_0_patreon]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_update_0_patreon
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_update_1_patreon]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_update_1_patreon
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

---------------------- Server 7 - Auto Updating - Friends ----------------------------------------------------------------

[program:character_friends_update_0_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_friends_update_0_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_friends_update_1_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_friends_update_1_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_friends_update_0_patreon]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_friends_update_0_patreon
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_friends_update_1_patreon]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_friends_update_1_patreon
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

---------------------- Server 8 - Auto Updating - Achievements ----------------------------------------------------------------

[program:character_achievements_update_0_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_achievements_update_0_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_achievements_update_1_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_achievements_update_1_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_achievements_update_0_patreon]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_achievements_update_0_patreon
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_achievements_update_1_patreon]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_achievements_update_1_patreon
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

---------------------- Server 9 - Auto Updating - Free Companies ----------------------------------------------------------------

[program:free_company_update_0_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest free_company_update_0_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:free_company_update_1_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest free_company_update_1_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:free_company_update_0_patron]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest free_company_update_0_patron
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:free_company_update_1_patron]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest free_company_update_1_patron
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

---------------------- Server 10 - Auto Updating - Linkshells ----------------------------------------------------------------

[program:linkshell_update_0_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest linkshell_update_0_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:linkshell_update_1_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest linkshell_update_1_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:linkshell_update_0_patron]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest linkshell_update_0_patron
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:linkshell_update_1_patron]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest linkshell_update_1_patron
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

---------------------- Server 11 - Auto Updating - PVP Team ----------------------------------------------------------------

[program:pvp_team_update_0_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest pvp_team_update_0_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:pvp_team_update_1_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest pvp_team_update_1_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:pvp_team_update_0_patron]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest pvp_team_update_0_patron
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:pvp_team_update_1_patron]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest pvp_team_update_1_patron
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log


---------------------- Server 12 - Auto Updating - Achievements ----------------------------------------------------------------

[program:character_achievements_update_2_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_achievements_update_2_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_achievements_update_3_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_achievements_update_3_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_achievements_update_4_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_achievements_update_4_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_achievements_update_5_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_achievements_update_5_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log
