# vm/SV_AutoUpdateAutoFriends.ini
---------------------- AutoUpdateAutoFriends | 202.182.123.22 -------------------------------------------------------

[program:character_update_0_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_update_0_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_update_1_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_update_1_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_update_2_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_update_2_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_update_3_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_update_3_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_update_4_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_update_4_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_update_5_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_update_5_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_friends_update_0_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_friends_update_0_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_friends_update_1_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_friends_update_1_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log
