# vm/SV_AutoUpdateFCLSPVP.ini
---------------------- AutoUpdateFCLSPVP | 45.76.199.186 -----------------------------------------------------------------

[program:pvp_team_update_0_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest pvp_team_update_0_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:pvp_team_update_1_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest pvp_team_update_1_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:free_company_update_0_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest free_company_update_0_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:free_company_update_1_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest free_company_update_1_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:linkshell_update_0_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest linkshell_update_0_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:linkshell_update_1_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest linkshell_update_1_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log
