# vm/SV_AutoUpdateAchievements.ini
---------------------- AutoUpdateAchievements | 45.76.107.66 -----------------------------------------------------------

[program:character_achievements_update_0_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_achievements_update_0_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_achievements_update_1_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_achievements_update_1_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_achievements_update_2_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_achievements_update_2_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_achievements_update_3_normal]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_achievements_update_3_normal
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log
