# vm/SV_Patreon.ini
---------------------- Patreon | 198.13.61.206 -------------------------------------------------------------

[program:character_friends_update_0_patreon]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_friends_update_0_patreon
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_friends_update_1_patreon]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_friends_update_1_patreon
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:pvp_team_update_0_patron]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest pvp_team_update_0_patron
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:pvp_team_update_1_patron]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest pvp_team_update_1_patron
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_update_0_patreon]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_update_0_patreon
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_update_1_patreon]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_update_1_patreon
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:free_company_update_0_patron]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest free_company_update_0_patron
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:free_company_update_1_patron]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest free_company_update_1_patron
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:linkshell_update_0_patron]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest linkshell_update_0_patron
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:linkshell_update_1_patron]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest linkshell_update_1_patron
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_achievements_update_0_patron]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_achievements_update_0_patron
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

[program:character_achievements_update_1_patron]
command=/usr/bin/php /home/dalamud/xivapi.com/bin/console AutoManagerRequest character_achievements_update_1_patron
user=dalamud
stderr_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_error.log
stdout_logfile=/home/dalamud/xivapi.com/AutoManagerRequest_output.log

